# scripts/crontab
*/7 * * * * cd $HOME/mor1.github.io-deployment && git pull
